# deploy/wifi-radar.service
[Unit]
Description=WiFi radar backend (ESP32 sensors + FastAPI on :8000)
Documentation=https://github.com/BugQuest/wifi-radar
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=realitynauts
Group=realitynauts
WorkingDirectory=/home/realitynauts/wifi-radar/backend
Environment=RADAR_BAUD=115200
# Static frontend lives in backend/static (built locally by `npm run build` and
# copied or symlinked there).  Data (DuckDB, parquet, config.json) stays in
# ~/wifi-radar-data/ via the RADAR_DATA env default — survives a repo wipe.
Environment=RADAR_STATIC=/home/realitynauts/wifi-radar/backend/static
# esptool for the in-UI firmware updater (panel 🔥).  We point at the
# pip-installed wrapper inside the IDF Python venv — it's executable, has a
# correct shebang to the venv interpreter, and brings esptool's deps along.
# The one at .../components/esptool_py/esptool/esptool.py is NOT +x.
# Adjust the version-tagged path (idf5.5_py3.14_env) if your IDF install differs.
Environment=RADAR_ESPTOOL=/home/realitynauts/.espressif/python_env/idf5.5_py3.14_env/bin/esptool.py
ExecStart=/home/realitynauts/wifi-radar/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
# needs access to /dev/ttyUSB* — keep the user in the dialout group instead of
# adding device ACL exceptions here
SupplementaryGroups=dialout

[Install]
WantedBy=multi-user.target
